# systemd/zauto-backup.timer
[Unit]
Description=Daily Zlatibor state backup

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
